# centaur-node/upstart/centaur.upstart
# centaur

start on stopped finish

respawn

exec /usr/bin/centaur-node


